# docker-image/docker.d/dc_provision/template/supervisor/conf.d/named.conf
[program:named]
command=/usr/sbin/named -f -4 -u bind
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
