# docker-image/docker.d/dc_provision/template/supervisor/conf.d/ntpd.conf
[program:ntpd]
command=/usr/sbin/ntpd -p /var/run/ntpd.pid -n -g
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
