# docker-image/docker.d/dc_join/template/cron/cron.d/sysvol-replication.template
*/5 * * * * root /usr/bin/rsync -XAavz --delete-after --password-file=/etc/samba/rsync.secret rsync://sysvol-replication@${RSYNCD_server}/SysVol/ /var/lib/samba/sysvol/
